# contrib/upstart/txgroupldap.conf
description "Twisted Group LDAP Provisioner"
author "Carl <waldbiec@lafayette.edu>"

start on runlevel [2345]
stop on runlevel [016]

kill timeout 60
respawn
respawn limit 10 5
 
script
/usr/bin/sudo -u grouper /bin/bash <<START
cd /opt/pychangelogger/provisioner/ 
. ./pyenv/bin/activate 
export LD_LIBRARY_PATH=/usr/local/lib64
./twistd.sh -n --syslog --prefix txgroupldap --pidfile /var/run/txgroupprovisioner/txgroupldap.pid  provisioner -c /etc/grouper/txgroupldap.cfg
START
end script

