# contrib/upstart/txuserldap.conf
description "Twisted User LDAP Provisioner"
author "Carl <waldbiec@lafayette.edu>"

start on runlevel [2345]
stop on runlevel [016]

kill timeout 60
respawn
respawn limit 10 5
 
script
/usr/bin/sudo -u grouper /bin/bash <<START
cd /opt/pychangelogger/provisioner/ 
. ./pyenv/bin/activate 
export LD_LIBRARY_PATH=/usr/local/lib64
./twistd.sh -n --syslog --prefix txuserldap --pidfile /var/run/txgroupprovisioner/txuserldap.pid  provisioner -c /etc/grouper/txuserldap.cfg
START
end script
